# scripts/claude-workbench.service.template
[Unit]
Description=Claude Workbench — Terminal Manager
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKDIR%/backend
Environment=PATH=%VENV%/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-%WORKDIR%/.env
ExecStart=%VENV%/bin/python -m uvicorn main:app --host 0.0.0.0 --port %PORT%
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
